# deploy/systemd/ghost-control.service
[Unit]
Description=Ghost Control API
After=network.target

[Service]
WorkingDirectory=/opt/ghost-app/ghost-app
EnvironmentFile=/opt/ghost-app/ghost-app/.env
Environment=JWT_ALGORITHM=RS256
ExecStart=/opt/ghost-app/ghost-app/venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 8000 --log-level warning
Restart=always

[Install]
WantedBy=multi-user.target
